# templates/default/etc_init_greyhole_watcher.conf.erb
# Greyhole service watcher/notifier
#
# Warning!!
#  Any local changes may be overwritten. File is managed by Opscode Chef.
#

description "Greyhole watcher/notifier"
author "Jason Kulatunga"

start on started greyhole
stop on stopping greyhole

script
    #this script doesn't do anything, just waits for greyhole to start/stop.
    sleep infinity
end script

<% if @pushover[:appid] && @pushover[:userid] %>
post-start script
    curl -s \
  --form-string "token=<%= @pushover[:appid]%>" \
  --form-string "user=<%= @pushover[:userid]%>" \
  --form-string "message=Greyhole has started" \
  https://api.pushover.net/1/messages.json

end script

post-stop script
    curl -s \
  --form-string "token=<%= @pushover[:appid]%>" \
  --form-string "user=<%= @pushover[:userid]%>" \
  --form-string "message=Greyhole has stopped" \
  https://api.pushover.net/1/messages.json

end script
<% end %>