# templates/default/etc_init_mount_shares_locally.conf.erb
# Mount Samba Shares Locally
#
# Warning!!
#  Any local changes may be overwritten. File is managed by Opscode Chef.
#
# based off of https://github.com/gboudreau/Greyhole/wiki/MountSharesLocally

description "Mount Samba Shares Locally"
author "Jason Kulatunga"

start on filesystem and static-network-up and started greyhole
stop on (runlevel [016] or stopping greyhole)

env RUN_AS=<%= @depot[:user] %>

pre-start script
    echo "Creating /mnt/samba folder if missing"
    if [ ! -x /mnt/samba/ ]; then
            mkdir -p /mnt/samba/
    fi
    <% if @pushover[:appid] && @pushover[:userid] %>
        curl -s \
        --form-string "token=<%= @pushover[:appid]%>" \
        --form-string "user=<%= @pushover[:userid]%>" \
        --form-string "message=Shares have been mounted locally" \
        https://api.pushover.net/1/messages.json
    <% end %>
    echo "finished prestart"
end script

script
    if [ -f /lib/lsb/init-functions ]; then
        . /lib/lsb/init-functions
    fi
    uid=`id -u $RUN_AS`
    gid=`id -g $RUN_AS`


    cd /mnt/samba/
    testparm -s /etc/samba/smb.conf 2>/dev/null | grep "^\[" | grep -v "\[global\]" | grep -v "\[homes\]" | awk -F'[' '{print $2}' | awk -F']' '{print $1}' | xargs -d "\n" mkdir -p
    sleep 5
    echo "Mounting Samba shares locally..."
    ls -1 | while read d; do
            if [ "`mount | grep "//127.0.0.1/$d/* on " | wc -l`" = "0" ]; then
                    /sbin/mount.cifs "//127.0.0.1/$d" "$d" --verbose -o credentials=/home/${RUN_AS}/.smb_credentials,uid=${uid},gid=${gid},file_mode=0660,dir_mode=0770,nobrl,hard,_netdev,iocharset=utf8,noserverino,mfsymlinks
            else
                    echo "  Share [$d] is already mounted."
            fi
    done
    echo "begin sleep cycle"
    sleep infinity
end script

post-stop script
 <% if @pushover[:appid] && @pushover[:userid] %>
    curl -s \
  --form-string "token=<%= @pushover[:appid]%>" \
  --form-string "user=<%= @pushover[:userid]%>" \
  --form-string "message=Shares have been unmounted" \
  https://api.pushover.net/1/messages.json
<% end %>

    echo "Unmounting locally mounted Samba shares..."
    /bin/umount -l /mnt/samba/*
    return 0
end script