# templates/default/etc_init_x11vnc.conf.erb
start on login-session-start

script

#this works in the vm.
#/usr/bin/x11vnc -xkb -forever -auth /var/run/lightdm/root/:1 -display :1 -rfbauth /etc/x11vnc.pass -rfbport 5900 -bg -o /var/log/x11vnc.log

#enable flags in alphabetical order (there's alot of them)
# -auth = The Xauthority file for the VNC server
# -display = X11 server display to connect to, usually :0
# -forever = means even if you disconnect from your session, x11vnc will still be alive
# -o = Write stderr messages to file "logfile" instead of to the terminal.
# -repeat = is allowing keyboard repeating keys (most annoying thing when it doesn’t work)
# -reopen = If the X server connection is disconnected, try to reopen the X display (up to one time.)
# -rfbauth = use authentication on RFB protocol
# -rfbport = TCP port for RFB protocol
# -xkb = will allow better clipboard interactions

#/usr/bin/x11vnc \
#-auth /var/run/lightdm/root/:1 \
#-display :1 \
#-forever \
#-o /var/log/x11vnc.log \
#-repeat \
#-rfbauth /etc/x11vnc.pass \
#-rfbport 5900 \
#-xkb

#this is an attempt to findauth when noone has logged in yet
/usr/bin/x11vnc \
-env FD_XDM=1 -findauth \
-display :1 \
-forever \
-o /var/log/x11vnc.log \
-repeat \
-rfbauth /etc/x11vnc.pass \
-rfbport 5900 \
-xkb



#/usr/bin/x11vnc -xkb -forever -auth /var/run/lightdm/root/:0 -display :0 -rfbauth /etc/x11vnc.pass -rfbport 5900 -bg -o /var/log/x11vnc.log
#x11vnc -noxdamage -many -display :0 -auth guess -forever -bg -rfbauth /etc/x11vnc.pass -rfbport 5900 -o ~/x11vnc.log

end script